# 02 Source/home/TMAP/BW/release/crontab/bw0.crontab.file
10 0 * * * /home/TMAP/BW/release/bin/FileCleanUp.sh TMAP-EM
10 0 1 * * /home/TMAP/BW/release/bin/TableCleanUP.sh
